# .local/share/scripts/vpn/portforward.service
[Unit]
Description=Port forwarding for vpn and qBittorrent

[Service]
Type=exec
ExecStart=/usr/local/bin/portforward.sh
Restart=always

[Install]
WantedBy=multi-user.target
