# .config/scripts/system/usr/local/lib/systemd/system/suwayomi-server.service
[Unit]
Description=A free and open source manga reader server that runs extensions built for Tachiyomi.
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=jellyfin
Group=jellyfin
ExecStart=/usr/bin/java -jar /usr/local/share/java/suwayomi-server/Suwayomi-Server.jar
Restart=on-failure

[Install]
WantedBy=multi-user.target
